# com.solancelab.claude-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd agent to keep the dashboard always running.
  EDIT the two paths below to match your machine:
    1. the absolute path to your `node` binary  (find it with: which node)
    2. the absolute path to this repo's server.mjs
  Then:
    cp com.solancelab.claude-dashboard.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.solancelab.claude-dashboard.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.solancelab.claude-dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/ABSOLUTE/PATH/TO/claude-config-dashboard/server.mjs</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>PORT</key>
    <string>7373</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
